# cmux/com.user.cmux-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.cmux-sync</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>__REPO_PATH__/cmux/sync-to-cmux.zsh</string>
	</array>
	<key>StartInterval</key>
	<integer>30</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/cmux-sync.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/cmux-sync.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
